# scripts/supervisord.conf
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=INFO
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[program:python-daemon]
command=/on_stop.py
autorestart=unexpected
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:mc-server]
command=/start
autorestart=unexpected
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
